# deploy/ec2/tokenization-api.service
[Unit]
Description=InlineRe-Identify Tokenization REST API
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Restart=always
RestartSec=5
# --network host (rather than -p 4051:4051) so PORT in the env file is the only place the port
# is configured, and so the container reaches the instance metadata service (for KMS instance-role
# credentials) as a normal host process would -- bridged containers need an extra IMDS hop-limit
# tweak at instance-launch time that this sidesteps entirely.
ExecStartPre=-/usr/bin/docker rm -f tokenization-api
ExecStart=/usr/bin/docker run --rm --name tokenization-api \
  --env-file /etc/inline-reidentify/tokenization-api.env \
  --network host \
  inline-reidentify-tokenization-api:latest
ExecStop=/usr/bin/docker stop -t 10 tokenization-api

[Install]
WantedBy=multi-user.target
